# service_checker.service
[Unit]
After=network.target

[Service]
ExecStart=/usr/local/bin/service_checker.sh

[Install]
WantedBy=default.target
